# deployment/reverb-worker.conf
; ==========================================================================
; Supervisor Configuration for Laravel Reverb WebSocket Server
; ==========================================================================
;
; This configuration runs the Laravel Reverb WebSocket server as a daemon
; process managed by Supervisor. Reverb provides real-time WebSocket
; functionality for broadcasting events to connected clients.
;
; Installation:
; 1. Install Supervisor: apt-get install supervisor (Debian/Ubuntu)
; 2. Copy this file to /etc/supervisor/conf.d/reverb-worker.conf
; 3. Update the paths below to match your deployment environment
; 4. Reload Supervisor: supervisorctl reread && supervisorctl update
; 5. Start the worker: supervisorctl start rackaudit-reverb
;
; Verify it's running:
;   supervisorctl status rackaudit-reverb
;
; View logs:
;   tail -f /var/log/rackaudit/reverb.log
;
; ==========================================================================

[program:rackaudit-reverb]
process_name=%(program_name)s
command=php /var/www/rackaudit/artisan reverb:start --host=0.0.0.0 --port=8080
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/rackaudit/reverb.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopwaitsecs=3600

; ==========================================================================
; Notes:
; ==========================================================================
;
; - Adjust the 'command' path to your application directory
; - Adjust 'user' to the web server user (www-data, nginx, apache, etc.)
; - Adjust 'stdout_logfile' path as needed
; - For production with TLS, configure a reverse proxy (Nginx/Caddy) to
;   handle SSL termination and proxy WebSocket connections to Reverb
;
; Example Nginx configuration for WebSocket proxy:
;
;   location /app {
;       proxy_pass http://127.0.0.1:8080;
;       proxy_http_version 1.1;
;       proxy_set_header Upgrade $http_upgrade;
;       proxy_set_header Connection "upgrade";
;       proxy_set_header Host $host;
;       proxy_set_header X-Real-IP $remote_addr;
;       proxy_set_header X-Forwarded-For $proxy_add_x_forwarded_for;
;       proxy_set_header X-Forwarded-Proto $scheme;
;       proxy_read_timeout 86400;
;   }
;
; ==========================================================================
